# roles/supervisor/files/api.conf
[program:api]
command=/vagrant/start.sh
directory=/vagrant
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/api.err.log
stdout_logfile=/var/log/supervisor/api.out.log
user=root
